# lib/udev/rules.d/99-net_up.rules
ACTION=="add|change", SUBSYSTEM=="net", ATTR{operstate}=="down", RUN+="upNet %k"
