# etc/xdg/autostart/netTun.desktop
[Desktop Entry]
Name=netTun
Exec=netTun
Type=Application
